# contrib/qalc_refresh_rates.service
[Unit]
Description="Refresh qalc exchange rates"

[Service]
ExecStart=/usr/bin/qalc -e
